# mattermost-bot.service
[Unit]
Description=Mattermost Bot for Agent Zero
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/a0/usr/workdir/mattermost-bot
Environment=MATTERMOST_URL=mattermost.house.svc.cluster.local
Environment=MATTERMOST_PORT=8065
Environment=MATTERMOST_TOKEN=YOUR_MATTERMOST_TOKEN_HERE
Environment=MATTERMOST_SCHEME=http
Environment=AGENT_ZERO_URL=http://agent-zero.house.svc.cluster.local
Environment=AGENT_ZERO_API_KEY=YOUR_AGENT_ZERO_API_KEY_HERE
Environment=LOG_LEVEL=INFO
Environment=PYTHONUNBUFFERED=1
ExecStart=/a0/usr/workdir/mattermost-bot/start.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
